# infra/launchd/dev.tommydoerr.influxdb-backup-r2-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.tommydoerr.influxdb-backup-r2-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Volumes/dev/influxdb/backup-r2-sync.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Volumes/dev/influxdb/backup-r2-sync.log</string>

  <key>StandardErrorPath</key>
  <string>/Volumes/dev/influxdb/backup-r2-sync.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
